# modules/PanCancerSystem/docker/supervisord.dev.conf
; =============================================================================
; OpenPanCan - Supervisor configuration (Development)
; =============================================================================
;
; Same as production but with hot-reload enabled for both frontend (Vite)
; and backend (uvicorn --reload).
;
; =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=debug

; ---------------------------------------------------------------------------
; Nginx - Web server (proxies to Vite dev server in dev mode)
; ---------------------------------------------------------------------------
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/var/log/supervisor/nginx-stdout.log
stderr_logfile=/var/log/supervisor/nginx-stderr.log
priority=1
stopasgroup=true
killasgroup=true

; ---------------------------------------------------------------------------
; Vite dev server - Frontend hot-reload
; ---------------------------------------------------------------------------
[program:vite]
command=npm run dev -- --host 0.0.0.0 --port 5173
directory=/app/frontend
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/var/log/supervisor/vite-stdout.log
stderr_logfile=/var/log/supervisor/vite-stderr.log
priority=2
stopasgroup=true
killasgroup=true

; ---------------------------------------------------------------------------
; Uvicorn - FastAPI backend with hot-reload
; ---------------------------------------------------------------------------
[program:uvicorn]
command=uvicorn backend.main:app --host 0.0.0.0 --port 8007 --reload --log-level debug
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/var/log/supervisor/uvicorn-stdout.log
stderr_logfile=/var/log/supervisor/uvicorn-stderr.log
priority=3
stopasgroup=true
killasgroup=true
environment=PYTHONUNBUFFERED="1",DATABASE_URL="sqlite:////data/openpancan.db"
