# modules/PanCancerSystem/docker/supervisord.conf
; =============================================================================
; OpenPanCan - Supervisor configuration
; =============================================================================
;
; Manages:
;   1. Nginx - Web server (frontend + API proxy)
;   2. Uvicorn - FastAPI backend server
;
; =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info

; ---------------------------------------------------------------------------
; Nginx - Web server
; ---------------------------------------------------------------------------
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/var/log/supervisor/nginx-stdout.log
stderr_logfile=/var/log/supervisor/nginx-stderr.log
priority=1
stopasgroup=true
killasgroup=true

; ---------------------------------------------------------------------------
; Uvicorn - FastAPI backend (PanCancerSystem)
; ---------------------------------------------------------------------------
[program:uvicorn]
command=uvicorn backend.main:app --host 0.0.0.0 --port 8007 --workers 4 --log-level info
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/var/log/supervisor/uvicorn-stdout.log
stderr_logfile=/var/log/supervisor/uvicorn-stderr.log
priority=2
stopasgroup=true
killasgroup=true
environment=PYTHONUNBUFFERED="1",DATABASE_URL="sqlite:////data/openpancan.db"

; ---------------------------------------------------------------------------
; Event listener for process state changes
; ---------------------------------------------------------------------------
[eventlistener:processes]
command=supervisor_stdout
buffer_size=100
events=PROCESS_STATE_STARTING,PROCESS_STATE_RUNNING,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
result_handler=supervisor_stdout:event_handler
